# example/smtpproxy.service
[Unit]
Description=SMTP Proxy

[Service]
EnvironmentFile=-/etc/default/smtpproxy
ExecStart=/usr/local/sbin/smtpproxy
User=nobody
