# example/smtpproxy.socket
[Unit]
Description=SMTP Proxy

[Socket]
ListenStream=25

[Install]
WantedBy=sockets.target
